# scripts/systemd/eq-log-maintain.timer
[Unit]
Description=Twice-weekly EverQuest log maintenance (for 99dps)

[Timer]
# Mondays and Thursdays at 11:00 local — twice a week. Rotation is size-gated
# (THRESHOLD), so a run only rotates a log that's actually oversized.
OnCalendar=Mon,Thu *-*-* 11:00:00
Persistent=true

[Install]
WantedBy=timers.target
